# cv/photonvision.service
# This file goes in /lib/systemd/system/photonvision.service
[Unit]
Description=PhotonVision

[Service]
Type=simple
ExecStart=/home/orangepi/.local/bin/uv run --project /app/cv /app/cv/server.py --master --instance "Instance 1" --port 8080 --config_file "/app/cv/photonvision_config.json" --ip_file "/app/cv/ip_suffix.txt" --field-layout k2026RebuiltWelded
Restart=on-failure

[Install]
WantedBy=multi-user.target
