# cv/photonvision2.service
# This file goes in /lib/systemd/system/photonvision2.service
[Unit]
Description=PhotonVision

[Service]
Type=simple
ExecStart=/home/orangepi/.local/bin/uv run --project /app/cv /app/cv/server.py --instance "Instance 2" --port 8081 --config_file "/app/cv/photonvision_config2.json" --field-layout k2026RebuiltWelded
Restart=on-failure

[Install]
WantedBy=multi-user.target
